# launchd/ai.openclaw.gateway-healthcheck.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.openclaw.gateway-healthcheck</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/openclaw/.openclaw/bin/openclaw_gateway_healthcheck.py</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/openclaw/.openclaw/logs/gateway-healthcheck.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/openclaw/.openclaw/logs/gateway-healthcheck.err.log</string>
  <key>WorkingDirectory</key>
  <string>/Users/openclaw/.openclaw</string>
</dict>
</plist>
